# internal/templates/scheduler.cron.tmpl
* * * * * deploy [ -f {{ .DeployPath }}/current/artisan ] && cd {{ .DeployPath }}/current && php artisan schedule:run >/dev/null 2>&1
